# deploy/supervisor-secuai.conf
; /etc/supervisor/conf.d/secuai-worker.conf
;
; Two queue workers tuned for a 2 GB VPS. Bump numprocs on bigger boxes.
;
; After installing:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start secuai-worker:*
;
; Why two pools?
;   - 'default' handles the bulk: notifications, audit log fan-out, webhooks.
;   - 'long' handles slow stuff: AI calls, cloud scans, data exports.
;     Separate pool prevents one slow scan from blocking 100 notifications.

[program:secuai-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/secuai/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --max-jobs=1000
autostart=true
autorestart=true
user=deploy
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/secuai-worker-default.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stopwaitsecs=3600

[program:secuai-worker-long]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/secuai/artisan queue:work redis --queue=ai,scans,exports --sleep=5 --tries=2 --max-time=3600 --memory=384
autostart=true
autorestart=true
user=deploy
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/secuai-worker-long.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stopwaitsecs=3600

[group:secuai-worker]
programs=secuai-worker-default,secuai-worker-long
